# deploy/kiosk-manager.service
[Unit]
Description=Kiosk Manager Dashboard
After=network.target

[Service]
Type=simple
User=kiosk
Group=kiosk
WorkingDirectory=/opt/kiosk-manager
ExecStart=/opt/kiosk-manager/venv/bin/gunicorn --bind 0.0.0.0:5000 --workers 1 --threads 4 --timeout 120 --access-logfile - --error-logfile - wsgi:app
Restart=always
RestartSec=5
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/kiosk/.Xauthority
Environment=XDG_RUNTIME_DIR=/run/user/1000
Environment=KIOSK_SECRET_KEY=change-me-to-random-string

[Install]
WantedBy=multi-user.target
